# packaging/debian/hub/postrm
#!/bin/sh
set -e

case "$1" in
    remove)
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
        fi
        ;;
    purge)
        if [ -d /run/systemd/system ]; then
            systemctl disable can-hub-web.service >/dev/null 2>&1 || true
            systemctl disable can-hub.service >/dev/null 2>&1 || true
            systemctl daemon-reload >/dev/null 2>&1 || true
        fi
        rm -rf /var/lib/can-hub
        ;;
esac

exit 0
