# packaging/debian/web/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if ! getent passwd can-hub >/dev/null 2>&1; then
            useradd --system --user-group --no-create-home \
                    --home-dir /var/lib/can-hub --shell /usr/sbin/nologin \
                    --comment "can-hub service" can-hub
        fi
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
            systemctl enable can-hub-web.service >/dev/null 2>&1 || true
            systemctl start can-hub-web.service >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
